# cliente/monitor-agent.service
[Unit]
Description=Demonio de heartbeat Monitor de Servidores
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=monitor-agent
Group=monitor-agent
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/monitor-agent/latidos.py
Restart=always
RestartSec=10
TimeoutStopSec=30

NoNewPrivileges=true
ProtectSystem=strict
# El agente escribe las claves Ed25519 en /etc/monitor-agent/ en el primer arranque.
# ReadWritePaths abre escritura solo en ese directorio manteniendo el resto protegido.
ReadWritePaths=/etc/monitor-agent
PrivateTmp=true
ProtectHome=true
MemoryDenyWriteExecute=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
